# data/protonplus.service
[Unit]
Description=ProtonPlus

[Service]
Type=oneshot
ExecStart={ExecStart}
