# data/protonplus.timer
[Unit]
Description=ProtonPlus

[Timer]
OnBootSec=1min
OnUnitActiveSec={OnUnitActiveSec}
Persistent=true

[Install]
WantedBy=timers.target
